# assets/_miwifi_r3l_all_23e37_2.8.51_INT.bin.extracted/squashfs-root/etc/uci-defaults/luci-i18n-korean
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ko_kr=한국
        commit luci
EOF

