# assets/_miwifi_r3l_firmware_5e1a6_2.13.7.bin.extracted/squashfs-root/etc/hotplug.d/iface/90-ustack
#!/bin/sh

if [ "$ACTION" = "ifup" -a "$INTERFACE" = "lan" ]; then
    # for http url
    [ -f /etc/init.d/ustack ] && /etc/init.d/ustack restart
fi
